# playbooks/roles/clustercheck/templates/etc/xinetd.d/clustercheck
# default: on
# description: mysql-chk
service mysql-chk
{
        disable = no
        flags = REUSE
        socket_type = stream
        port = 9200
        wait = no
        user = nobody
        server = /usr/local/sbin/mysql-chk
        log_on_failure += USERID
        only_from = 0.0.0.0/0
        per_source = UNLIMITED
}
